# apps/cli/src/templates/watchdog/ops/watchdog/tx-ralph-watchdog.service
[Unit]
Description=tx RALPH Watchdog
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=/bin/bash __PROJECT_DIR__/scripts/watchdog-launcher.sh start --foreground
ExecStop=/bin/bash __PROJECT_DIR__/scripts/watchdog-launcher.sh stop
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
